# debian/multiflexi-executor.service
#!/bin/sh
# If this exists, it is installed into lib/systemd/system/multiflexi.service in multiflexi
# See dh_systemd_enable(1).

[Unit]
Description=Run MultiFlexi scheduled jobs
After=multi-user.target

[Service]
User=multiflexi
Group=multiflexi
EnvironmentFile=/etc/multiflexi/multiflexi.env
ExecStart=/usr/bin/php /usr/lib/multiflexi-executor/daemon.php
LimitNOFILE=8192:16384
Type=simple
# Hard memory ceiling to prevent uncontrolled growth; adjust as needed.
MemoryMax=2G
# Propagate soft limit (in MB) to daemon for graceful early exit (matches daemon code usage)
Environment=MULTIFLEXI_MEMORY_LIMIT_MB=1800
# Automatically restart on exit (including after hitting memory limit)
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
